# _retired/ansible/roles/motorcade_postgres/templates/motorcade-postgres.container.j2
[Unit]
Description=Motorcade Postgres (official image)
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=motorcade-postgres
Image={{ legacy_pg_image }}
Pull=newer
PublishPort={{ postgres_bind_host }}:{{ postgres_port }}:5432
Environment=POSTGRES_PASSWORD={{ vault_postgres_password }}
Environment=POSTGRES_USER=motorcade
Environment=POSTGRES_DB=motorcade
Volume={{ motorcade_vol_dir }}/postgres:/var/lib/postgresql/data:rw

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
